# pkgbuilds/sdf/sdf.service
[Unit]
Description=Simple Dell Fan

[Service]
ExecStart=/usr/bin/sdf
Restart=on-failure

[Install]
WantedBy=multi-user.target
